# tools/echo_relay/echo-relay.service.example
[Unit]
Description=ECHO Secure Relay API (localhost-only, disabled by default)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/TITAN
Environment=ECHO_RELAY_ENABLED=false
Environment=ECHO_INTERNAL_BASE_URL=http://127.0.0.1:8765
Environment=ECHO_RELAY_API_KEY=
Environment=ECHO_INTERNAL_API_KEY=
ExecStart=/usr/bin/python -m uvicorn titan_echo.echo_relay_api:app --host 127.0.0.1 --port 8766
Restart=no
User=titan
Group=titan

[Install]
WantedBy=multi-user.target
